# config/systemd/rovac-edge-uros-agent.service
[Unit]
Description=ROVAC micro-ROS Agent (UDP bridge for wireless ESP32 nodes)
Wants=network-online.target
After=network-online.target
# Mux depends on Agent being up — Agent must start before mux
Before=rovac-edge-mux.service
PartOf=rovac-edge.target
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi

# DDS Configuration (Agent bridges XRCE-DDS to CycloneDDS)
Environment=ROS_DOMAIN_ID=42
Environment=RMW_IMPLEMENTATION=rmw_cyclonedds_cpp
Environment=CYCLONEDDS_URI=file:///home/pi/robots/rovac/config/cyclonedds_pi.xml

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rovac-uros-agent

# Timeout settings
TimeoutStartSec=15
TimeoutStopSec=10

# Resource limits for Pi 5
# Steady-state RSS is ~65-70MB (CycloneDDS discovery cache for all network
# participants). MemoryHigh is set well above baseline to avoid false pressure.
# The watchdog script provides a softer restart at 150MB before the hard kill.
MemoryHigh=200M
MemoryMax=250M
CPUQuota=20%

# Restart policy
Restart=always
RestartSec=3

# UDP port 8888 — must match ESP32 NVS config (agent_port)
ExecStart=/bin/bash -c 'source /opt/ros/jazzy/setup.bash && source /home/pi/uros_ws/install/setup.bash && exec /home/pi/uros_ws/install/micro_ros_agent/lib/micro_ros_agent/micro_ros_agent udp4 --port 8888 -v0'

# After Agent starts (or restarts), reset the ESP32 to force a fresh
# XRCE-DDS session. Without this, a restarted Agent has no session state
# and silently drops all ESP32 data until the ESP32's ping detects the loss.
ExecStartPost=/home/pi/robots/rovac/scripts/edge/reset_esp32_motor.sh

[Install]
WantedBy=rovac-edge.target
