# config/systemd/rovac-edge-map-tf.service
[Unit]
Description=ROVAC Map->Odom Static Transform (fallback when SLAM not running)
After=network.target rovac-edge-tf.service
PartOf=rovac-edge.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
Environment=ROS_DOMAIN_ID=42
Environment=RMW_IMPLEMENTATION=rmw_cyclonedds_cpp
Environment=CYCLONEDDS_URI=file:///home/pi/cyclonedds_pi.xml
ExecStart=/home/pi/launch_map_odom_tf.sh
Restart=always
RestartSec=2

[Install]
WantedBy=rovac-edge.target
