# config/systemd/rovac-edge-tf-relay.service
[Unit]
Description=ROVAC TF Relay (best_effort → reliable for ESP32 micro-ROS)
After=rovac-edge-uros-agent.service
PartOf=rovac-edge.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
Environment=ROS_DOMAIN_ID=42
Environment=RMW_IMPLEMENTATION=rmw_cyclonedds_cpp
Environment=CYCLONEDDS_URI=file:///home/pi/robots/rovac/config/cyclonedds_pi.xml
ExecStart=/bin/bash -c 'source /opt/ros/jazzy/setup.bash && exec python3 /home/pi/robots/rovac/scripts/edge/tf_relay.py'
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rovac-tf-relay

[Install]
WantedBy=rovac-edge.target
