# config/systemd/rovac-edge.target
[Unit]
Description=ROVAC Edge Stack (Pi 5 Sensors and Motors)
Documentation=https://github.com/rovac-robot
After=network-online.target
Wants=network-online.target

# Core services
Wants=rovac-edge-hiwonder.service
Wants=rovac-edge-lidar.service
Wants=rovac-edge-supersensor.service
Wants=rovac-edge-mux.service
Wants=rovac-edge-tf.service
Wants=rovac-edge-obstacle.service
Wants=rovac-edge-stereo-depth.service
Wants=rovac-edge-stereo-obstacle.service

# Phone sensors (optional - starts if phone connected)
Wants=rovac-edge-phone-sensors.service

# Camera (optional)
Wants=rovac-camera.service

[Install]
WantedBy=multi-user.target
