# config/systemd/rovac-camera.service
[Unit]
Description=ROVAC Phone Camera (scrcpy + v4l2loopback + ROS2 publisher)
Wants=network-online.target
After=network-online.target
PartOf=rovac-edge.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
PermissionsStartOnly=true

# Logging configuration
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rovac-camera

# Timeout settings
TimeoutStartSec=60
TimeoutStopSec=15

# Resource limits for Pi 5
MemoryMax=512M
CPUQuota=80%

# Restart policy
Restart=always
RestartSec=5

# Defaults (override via: systemctl edit rovac-camera.service)
Environment=PHONE_CAMERA_ID=0
Environment=ROVAC_VIDEO_DEVICE=/dev/video10

# Ensure v4l2loopback exists and is writable.
ExecStartPre=/bin/bash -lc 'modprobe -r v4l2loopback 2>/dev/null || true; modprobe v4l2loopback devices=1 video_nr=10 card_label="Phone_Camera" exclusive_caps=1; chmod 666 /dev/video10'

ExecStart=/bin/bash -lc '\
set -e; set -o pipefail; \
source /opt/ros/jazzy/setup.bash; \
[ -f /home/pi/ros2_env.sh ] && source /home/pi/ros2_env.sh; \
VIDEO_DEV="${ROVAC_VIDEO_DEVICE:-/dev/video10}"; \
CAM_ID="${PHONE_CAMERA_ID:-0}"; \
ADB="/usr/bin/adb"; \
SCRCPY="/usr/local/bin/scrcpy"; \
PUBLISHER="/home/pi/phone_camera_publisher.py"; \
\
$ADB start-server >/dev/null 2>&1 || true; \
phone_ok=0; \
for _ in $(seq 1 30); do \
  if $ADB devices 2>/dev/null | awk '"'"'NR>1 {print $2}'"'"' | grep -q "^device$"; then \
    phone_ok=1; \
    break; \
  fi; \
  echo "[rovac-camera] Waiting for phone via ADB..."; \
  sleep 2; \
done; \
if [ "$phone_ok" != "1" ]; then \
  echo "[rovac-camera] No phone detected via ADB; exiting"; \
  exit 1; \
fi; \
if [ ! -x "$SCRCPY" ]; then \
  echo "[rovac-camera] scrcpy not found at $SCRCPY"; \
  exit 1; \
fi; \
if [ ! -f "$PUBLISHER" ]; then \
  echo "[rovac-camera] missing publisher: $PUBLISHER"; \
  exit 1; \
fi; \
if [ ! -e "$VIDEO_DEV" ]; then \
  echo "[rovac-camera] missing video device: $VIDEO_DEV"; \
  exit 1; \
fi; \
\
$SCRCPY --video-source=camera --camera-id="$CAM_ID" --camera-size=640x480 --camera-fps=15 --v4l2-sink="$VIDEO_DEV" --no-window --no-audio >>/tmp/phone_cam.log 2>&1 & \
SCRCPY_PID=$!; \
cleanup() { \
  if kill -0 "$SCRCPY_PID" 2>/dev/null; then \
    kill "$SCRCPY_PID" 2>/dev/null || true; \
    sleep 0.5; \
    kill -9 "$SCRCPY_PID" 2>/dev/null || true; \
  fi; \
}; \
trap cleanup EXIT INT TERM; \
exec /usr/bin/python3 "$PUBLISHER" "$VIDEO_DEV" "/phone/camera/image_raw" >>/tmp/camera_pub.log 2>&1'

[Install]
WantedBy=rovac-edge.target
