# config/systemd/rovac-edge-uros-agent.service
[Unit]
Description=ROVAC micro-ROS Agent (UDP bridge for ESP32 motor wireless)
Wants=network-online.target
After=network-online.target
PartOf=rovac-edge.target
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi

# DDS Configuration (Agent bridges XRCE-DDS to CycloneDDS)
Environment=ROS_DOMAIN_ID=42
Environment=RMW_IMPLEMENTATION=rmw_cyclonedds_cpp
Environment=CYCLONEDDS_URI=file:///home/pi/robots/rovac/config/cyclonedds_pi.xml

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rovac-uros-agent

# Timeout settings
TimeoutStartSec=15
TimeoutStopSec=10

# Resource limits for Pi 5
MemoryMax=128M
CPUQuota=20%

# Restart policy
Restart=always
RestartSec=3

# UDP port 8888 — must match ESP32 NVS config (agent_port)
ExecStart=/bin/bash -lc 'source /opt/ros/jazzy/setup.bash; source /home/pi/uros_ws/install/setup.bash; exec micro_ros_agent udp4 --port 8888 -v4'

[Install]
WantedBy=rovac-edge.target
