# config/systemd/rovac-edge-esp32.service
[Unit]
Description=ROVAC ESP32 AT8236 Motor Driver (USB-serial motor + encoder)
After=network.target dev-esp32_motor.device
BindsTo=dev-esp32_motor.device
PartOf=rovac-edge.target
Before=rovac-edge-mux.service rovac-edge-obstacle.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/robots/rovac/hardware/esp32_at8236_driver

Environment=ROS_DOMAIN_ID=42
Environment=RMW_IMPLEMENTATION=rmw_cyclonedds_cpp
Environment=CYCLONEDDS_URI=file:///home/pi/robots/rovac/config/cyclonedds_pi.xml

ExecStart=/home/pi/robots/rovac/hardware/esp32_at8236_driver/launch_esp32_at8236.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rovac-esp32

MemoryMax=256M
CPUQuota=50%

[Install]
WantedBy=rovac-edge.target
