# config/systemd/rovac-edge-supersensor.service
[Unit]
Description=ROVAC Super Sensor Node (4x ultrasonic, LED, servo)
Wants=network-online.target
After=network-online.target
PartOf=rovac-edge.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/robots/rovac/hardware/super_sensor
Environment=ROS_DOMAIN_ID=42
Environment=RMW_IMPLEMENTATION=rmw_cyclonedds_cpp
Environment=CYCLONEDDS_URI=file:///home/pi/robots/rovac/config/cyclonedds_pi.xml
ExecStart=/bin/bash -lc 'source /opt/ros/jazzy/setup.bash; [ -f /home/pi/robots/rovac/config/ros2_env.sh ] && source /home/pi/robots/rovac/config/ros2_env.sh; exec python3 /home/pi/robots/rovac/hardware/super_sensor/super_sensor_ros2_node.py --ros-args -p port:=/dev/ttyUSB1'
Restart=always
RestartSec=3

[Install]
WantedBy=rovac-edge.target
