# config/systemd/rovac-edge-ps2-joy.service
[Unit]
Description=ROVAC PS2 Controller joy_node (USB receiver -> /joy)
Documentation=https://github.com/mohammednazmy/rovac
Wants=network-online.target
After=network-online.target
Before=rovac-edge-ps2-mapper.service
PartOf=rovac-edge.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi

# DDS Configuration
Environment=ROS_DOMAIN_ID=42
Environment=RMW_IMPLEMENTATION=rmw_cyclonedds_cpp
Environment=CYCLONEDDS_URI=file:///home/pi/robots/rovac/config/cyclonedds_pi.xml

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rovac-ps2-joy

# Resource limits
MemoryMax=64M
CPUQuota=10%

Restart=always
RestartSec=3

ExecStart=/bin/bash -lc 'source /opt/ros/jazzy/setup.bash; exec ros2 run joy joy_node --ros-args -p device_name:=/dev/input/js0 -p deadzone:=0.05 -p autorepeat_rate:=20.0 -p coalesce_interval_ms:=2'

[Install]
WantedBy=rovac-edge.target
