# robot_mcp_server/robot-mcp.service
[Unit]
Description=Robot MCP Server for ChatGPT Voice Control
After=network.target ros2.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/robots/rovac/robot_mcp_server
ExecStart=/home/pi/robots/rovac/robot_mcp_server/start_mcp_server.sh
Restart=on-failure
RestartSec=5
Environment=HOME=/home/pi

[Install]
WantedBy=multi-user.target
